# crontab
0 18 * * * cd /app && npm run daily-sync